# launchd/com.stablemischief.relationship-intel.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stablemischief.relationship-intel.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mitchdornich/Documents/GitHub/SM-CRM-System/scripts/relationship-intel-daily.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mitchdornich/Documents/GitHub/SM-CRM-System</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mitchdornich/Documents/GitHub/SM-CRM-System/output/logs/relationship-intel-daily.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mitchdornich/Documents/GitHub/SM-CRM-System/output/logs/relationship-intel-daily.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
